# config/desktop-watcher.service
[Unit]
Description=DesktopWatcher activity tracker
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/desktop-watcherd
Restart=on-failure
RestartSec=5

# Give the daemon access to the X display
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=graphical-session.target
